# deploy/systemd/studio-bei-bitable-sync.timer
[Unit]
Description=Studio Bei OS · 每 5 分钟同步一次飞书多维表格

[Timer]
# 启动后 1 min 触发首次，之后每 5 min
OnBootSec=1min
OnUnitActiveSec=5min
Unit=studio-bei-bitable-sync.service
Persistent=true

[Install]
WantedBy=timers.target
